# dpkg/prototype/DEBIAN/postinst
#!/bin/sh
update-rc.d saltnode-setup defaults

