# dpkg/prototype/DEBIAN/prerm
#!/bin/sh
update-rc.d -f saltnode-setup remove

